# rules/macos/persistence_startup_item_plist_creation.toml
[metadata]
creation_date = "2026/01/30"
integration = ["endpoint"]
maturity = "production"
updated_date = "2026/01/30"

[rule]
author = ["Elastic"]
description = """
Detects the creation or modification of a StartupParameters.plist file, indicating the presence of a 
StartupItem on the system. StartupItems have been deprecated on modern macOS systems (post Mavericks) 
in favor of Launch Daemons but still function. Creation of a StartupItem should be highly suspicious 
as legitimate applications no longer use this method for persistence.
"""
from = "now-9m"
index = ["logs-endpoint.events.file-*"]
language = "eql"
license = "Elastic License v2"
name = "Suspicious StartupItem Plist Creation"
references = [
    "https://www.virusbulletin.com/uploads/pdf/conference/vb2014/VB2014-Wardle.pdf"
]
risk_score = 73
rule_id = "15606250-449d-46a8-aaff-4043e42aefb9"
severity = "high"
tags = [
    "Domain: Endpoint",
    "OS: macOS",
    "Use Case: Threat Detection",
    "Tactic: Persistence",
    "Data Source: Elastic Defend"
]
timestamp_override = "event.ingested"
type = "eql"
note = """## Triage and analysis

> **Disclaimer**:
> This investigation guide was created using generative AI technology and has been reviewed to improve its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit your specific environment and operational needs.

### Investigating Suspicious StartupItem Plist Creation

This rule detects the creation of StartupItem plist files, which is a deprecated but still functional persistence mechanism on macOS. Threat actors may use this technique to establish persistence as it is less commonly monitored than LaunchAgents.

### Possible investigation steps

- Examine the StartupItem plist contents and associated scripts.
- Identify the binary or command configured to execute at startup.
- Review the process that created the StartupItem.
- Check for other persistence mechanisms on the system.

### False positive analysis

- Legacy applications may still use StartupItems for compatibility.
- Some enterprise software uses deprecated startup mechanisms.

### Response and remediation

- Remove the malicious StartupItem plist and associated files.
- Reboot to verify the persistence was successfully removed.
- Investigate the initial access vector.
- Monitor for recreation of the StartupItem.
"""
query = '''
file where host.os.type == "macos" and event.type != "deletion" and
  file.path like ("/System/Library/StartupItems/*/StartupParameters.plist",
                  "/Library/StartupItems/*/StartupParameters.plist") and
  not process.code_signature.signing_id == "com.apple.shove"
'''

[[rule.threat]]
framework = "MITRE ATT&CK"

  [rule.threat.tactic]
  name = "Persistence"
  id = "TA0003"
  reference = "https://attack.mitre.org/tactics/TA0003/"

  [[rule.threat.technique]]
  name = "Boot or Logon Initialization Scripts"
  id = "T1037"
  reference = "https://attack.mitre.org/techniques/T1037/"

    [[rule.threat.technique.subtechnique]]
    name = "Startup Items"
    id = "T1037.005"
    reference = "https://attack.mitre.org/techniques/T1037/005/"
