# rules/linux/persistence_rc_local_error_via_syslog.toml
[metadata]
creation_date = "2024/06/21"
integration = ["system"]
maturity = "production"
updated_date = "2025/01/08"

[rule]
author = ["Elastic"]
description = """
This rule monitors the syslog log file for error messages related to the rc.local process. The rc.local file is a script
that is executed during the boot process on Linux systems. Attackers may attempt to modify the rc.local file to execute
malicious commands or scripts during system startup. This rule detects error messages such as "Connection refused," "No
such file or directory," or "command not found" in the syslog log file, which may indicate that the rc.local file has
been tampered with.
"""
from = "now-9m"
index = ["logs-system.syslog-*"]
language = "kuery"
license = "Elastic License v2"
name = "Suspicious rc.local Error Message"
references = [
    "https://www.intezer.com/blog/malware-analysis/hiddenwasp-malware-targeting-linux-systems/",
    "https://pberba.github.io/security/2022/02/06/linux-threat-hunting-for-persistence-initialization-scripts-and-shell-configuration/#8-boot-or-logon-initialization-scripts-rc-scripts",
    "https://www.cyberciti.biz/faq/how-to-enable-rc-local-shell-script-on-systemd-while-booting-linux-system/",
    "https://www.elastic.co/security-labs/sequel-on-persistence-mechanisms",
]
risk_score = 21
rule_id = "69c116bb-d86f-48b0-857d-3648511a6cac"
setup = """## Setup

This rule requires data coming in from one of the following integrations:
- Filebeat

### Filebeat Setup
Filebeat is a lightweight shipper for forwarding and centralizing log data. Installed as an agent on your servers, Filebeat monitors the log files or locations that you specify, collects log events, and forwards them either to Elasticsearch or Logstash for indexing.

#### The following steps should be executed in order to add the Filebeat for the Linux System:
- Elastic provides repositories available for APT and YUM-based distributions. Note that we provide binary packages, but no source packages.
- To install the APT and YUM repositories follow the setup instructions in this [helper guide](https://www.elastic.co/guide/en/beats/filebeat/current/setup-repositories.html).
- To run Filebeat on Docker follow the setup instructions in the [helper guide](https://www.elastic.co/guide/en/beats/filebeat/current/running-on-docker.html).
- To run Filebeat on Kubernetes follow the setup instructions in the [helper guide](https://www.elastic.co/guide/en/beats/filebeat/current/running-on-kubernetes.html).
- For quick start information for Filebeat refer to the [helper guide](https://www.elastic.co/guide/en/beats/filebeat/8.11/filebeat-installation-configuration.html).
- For complete Setup and Run Filebeat information refer to the [helper guide](https://www.elastic.co/guide/en/beats/filebeat/current/setting-up-and-running.html).

#### Rule Specific Setup Note
- This rule requires the Filebeat System Module to be enabled.
- The system module collects and parses logs created by the system logging service of common Unix/Linux based distributions.
- To run the system module of Filebeat on Linux follow the setup instructions in the [helper guide](https://www.elastic.co/guide/en/beats/filebeat/current/filebeat-module-system.html).
"""
severity = "low"
tags = ["Domain: Endpoint", "OS: Linux", "Use Case: Threat Detection", "Tactic: Persistence"]
timestamp_override = "event.ingested"
type = "query"

query = '''
host.os.type:linux and event.dataset:system.syslog and process.name:rc.local and
message:("Connection refused" or "No such file or directory" or "command not found")
'''
note = """## Triage and analysis

### Disclaimer

This investigation guide was generated using generative AI technology and has been reviewed to improve its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit your specific environment and operational needs.

### Investigating Suspicious rc.local Error Message

The rc.local script is crucial in Linux systems, executing commands at boot. Adversaries may exploit this by inserting malicious scripts to gain persistence. The detection rule monitors syslog for specific error messages linked to rc.local, such as "Connection refused," indicating potential tampering. This helps identify unauthorized modifications, aligning with MITRE ATT&CK's persistence tactics.

### Possible investigation steps

- Review the syslog entries around the time of the alert to gather more context about the suspicious error message and identify any other related log entries.
- Verify the integrity of the rc.local file by comparing its current state with a known good backup or using a file integrity monitoring tool to check for unauthorized changes.
- Use Osquery to list the current contents of the rc.local file with the query: `SELECT * FROM file WHERE path='/etc/rc.local';` to identify any suspicious or unexpected entries.
- Investigate the process tree to identify any parent or child processes related to rc.local using the query fields `process.name:rc.local` and `host.os.type:linux` to trace the execution flow.
- Check for recent modifications to the rc.local file by examining the file's metadata, such as the last modified timestamp, using `ls -l /etc/rc.local` or an equivalent command.
- Search for any recent user account activity or privilege escalations around the time of the alert that could indicate unauthorized access or changes to the rc.local file.
- Analyze network logs for any unusual outbound connections or failed connection attempts that coincide with the "Connection refused" error message.
- Review the system's boot logs to identify any anomalies or errors during the startup process that could be related to the rc.local script execution.
- Investigate other system logs for any additional error messages like "No such file or directory" or "command not found" that might provide further context on the issue.
- Correlate the findings with known MITRE ATT&CK techniques, specifically focusing on persistence tactics, to determine if the activity aligns with known adversary behaviors.

### False positive analysis

- Routine administrative tasks or legitimate software updates may trigger error messages like "No such file or directory" if certain expected files are temporarily unavailable during the update process. Users can handle these by creating exceptions for known maintenance windows or specific update processes.
- Custom scripts or applications that are intentionally executed at boot but occasionally fail due to misconfigurations might generate "command not found" errors. Users should review these scripts and, if deemed non-threatening, exclude them from triggering alerts by specifying the script names or paths in the exception list.
- Network-related services that are not yet fully initialized during the boot process might cause "Connection refused" errors. Users can mitigate these false positives by adjusting the timing of service startups or excluding specific network services from the rule if they are known to cause such errors without security implications.

### Response and remediation

- Immediately isolate the affected system from the network to prevent further unauthorized access or spread of malicious activity.
- Conduct a thorough investigation of the rc.local file to identify any unauthorized modifications or suspicious entries.
- Review system logs, particularly syslog, to trace the origin of the suspicious activity and identify any other affected systems or files.
- Restore the rc.local file from a known good backup if available, ensuring that it is free from any unauthorized changes.
- Update and patch the system to the latest security standards to mitigate any known vulnerabilities that could have been exploited.
- Implement strict access controls and permissions for critical system files like rc.local to prevent unauthorized modifications.
- Enhance logging policies to ensure comprehensive monitoring of system changes and potential security incidents, integrating with SIEM solutions for real-time alerts.
- Conduct a security audit of the system to identify and address any other potential security weaknesses or misconfigurations.
- Escalate the incident to the security operations team or relevant authorities if the investigation reveals a broader threat or compliance breach.
- Educate and train staff on recognizing and responding to security incidents, emphasizing the importance of maintaining system integrity and vigilance against persistence tactics."""


[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1037"
name = "Boot or Logon Initialization Scripts"
reference = "https://attack.mitre.org/techniques/T1037/"
[[rule.threat.technique.subtechnique]]
id = "T1037.004"
name = "RC Scripts"
reference = "https://attack.mitre.org/techniques/T1037/004/"



[rule.threat.tactic]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/tactics/TA0003/"

